# roles/desktop/hyprland/files/home/.config/systemd/user/hyprpaper.service
[Unit]
Description=Hyprland random wallpaper switcher

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/hyprpaper-hourly.sh
